# configs/celery.conf
; This is an example Celery Configuration
; Move to /etc/supervisor/conf.d/celery.conf

[program:celery]
command=/path/to/virtualenv/bin/celery -A app.celery worker --loglevel=INFO

user=vagrant
environment=SNIPSPACE_SETTINGS=/path/to/snipspace_settings.py
directory=/path/to/snip.space
numprocs=1

; This directory and file has to be created
stdout_logfile=/var/log/celery/worker.log
stderr_logfile=/var/log/celery/worker.log

autostart=true
autorestart=true
startsecs=10

stopwaitsecs = 100

